# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true

[program:queue-worker]
command=php artisan queue:work
autostart=true
autorestart=true
